# rrat.service
[Unit]
Description=RRAT Demon

[Service]
ExecStart=%h/rratfiles/rrat-demon

[Install]
WantedBy=default.target